# rootfs/etc/gotpl/root.crontab.tmpl
SHELL=/bin/bash

*/10 * * * * /opt/bin/default_cert_create
0 0 * * 1 /opt/bin/default_cert_renew